# meta-quanta/meta-olympus-nuvoton/recipes-phosphor/sel-logger/phosphor-sel-logger/xyz.openbmc_project.Logging.IPMI.service
[Unit]
Description=IPMI SEL Logging Service
After=olympus-reload-sensor.service

[Service]
Restart=always
ExecStart=/usr/bin/sel-logger
Type=simple

[Install]
WantedBy=multi-user.target
