# meta-quanta/meta-olympus-nuvoton/recipes-olympus-nuvoton/network/olympus-nuvoton-iptable-save/phosphor-save-iptable-rules.service
[Unit]
Description=Save iptable rules
DefaultDependencies=no
Before=reboot.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/sh -c "iptables-save > /etc/iptables/iptables.rules"

[Install]
WantedBy=reboot.target
