# meta-quanta/meta-olympus-nuvoton/recipes-phosphor/fans/phosphor-pid-control/phosphor-pid-control.service
[Unit]
Description=OpenBMC Fan Control Daemon
After=obmc-chassis-poweron@0.target

[Service]
Type=simple
ExecStartPre=/usr/bin/env mapper wait /xyz/openbmc_project/sensors/temperature/Core_0_CPU0
ExecStart=/usr/bin/swampd
Restart=always
RestartSec=5
StartLimitInterval=0
ExecStopPost=/usr/bin/fan-default-speed.sh

