# meta-evb/meta-evb-nuvoton/meta-buv-runbmc/recipes-phosphor/sel-logger/phosphor-sel-logger/xyz.openbmc_project.Logging.IPMI.service
[Unit]
Description=IPMI SEL Logging Service
Wants=rsyslog.service
After=rsyslog.service

[Service]
Restart=always
ExecStart=/usr/bin/sel-logger
Type=simple

[Install]
WantedBy=multi-user.target
